# etc/systemd/system/e2pv.service
[Unit]
After=network.target

[Service]
ExecStart=/home/pi/enecsys/e2pv.sh

[Install]
WantedBy=default.target
